# omni_platform_from_E/omni-auto-launcher.service
[Unit]
Description=OMNI Platform Auto Launcher - Self-Healing AI Platform
After=network.target
Wants=network.target

[Service]
Type=simple
User=omni
Group=omni
WorkingDirectory=/opt/omni_platform
Environment=PATH=/opt/omni_platform/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONPATH=/opt/omni_platform
ExecStart=/opt/omni_platform/venv/bin/python /opt/omni_platform/omni_auto_launcher.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omni-auto-launcher

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/omni_platform /tmp /var/log/omni_platform

# Resource limits
MemoryLimit=2G
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
